# etc/cron.daily/fstrim
#!/bin/sh
fstrim -v / >> /var/log/fstrim.log
echo `date` 'TRIM has executed' >> /var/log/fstrim.log
