# etc/supervisor/backgrounds.conf
[program:backgrounds]
command=/home/kick/Documents/dotfiles/etc/xinitrc.d/52-background-image.sh
user = kick
environment = DISPLAY=":0",HOME=/home/kick
