# etc/supervisor/laptop-polybar.conf
[program:laptop-polybar]
command = polybar laptop
user = kick
environment = DISPLAY=":0",HOME=/home/kick
