# etc/supervisor/main-polybar.conf
[program:main-polybar]
command = polybar main
user = kick
environment = DISPLAY=":0",HOME=/home/kick
