# etc/supervisor/unclutter.conf
[program:unclutter]
command=/usr/bin/unclutter
user = kick
environment = DISPLAY=":0",HOME=/home/kick
